# supervisor.d/rabbitmq.conf
[program:rabbitmq]
command=/usr/sbin/rabbitmq-server 
autostart=true
autorestart=true
startsecs=1
startretries=3

